# docker/supervisor.d/workers.ini
[program:mog_worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php71 /var/www/fsm/public_html/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/fsm/logs/mog_worker.log

[program:mog_test_worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php71 /var/www/fsm-test/public_html/artisan queue:work redis --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/fsm-test/logs/mog_worker.log

